# .config/systemd/user/dotfiles_backup.service
[Service]
# Type = simple
ExecStart=%h/my_scripts/usw/shell/dotfiles_backup.sh

; not required since the timer is enabled now
; [Install]
; WantedBy = default.target
; WantedBy = multi-user.target
