# .config/systemd/user/email_notify.timer
[Unit]
After = graphical.target network.target

[Timer]
# Every 15 minutes
OnCalendar = *:0/15

[Install]
WantedBy = timers.target
