# .config/systemd/user/email_notify.service
[Service]
Environment=DISPLAY=:0
ExecStart = %h/.local/bin/email_notify
