# .config/systemd/user/fail@.service
[Service]
ExecStart = %h/.local/bin/notify '%f.service has failed' 0 10000
