# .config/systemd/user/pacman_updates_notify.service
[Service]
ExecStart=%h/.local/bin/sh/pacman_updates_notify.sh
