# .config/systemd/user/space_used.service
[Service]
ExecStart=%h/.local/bin/sh/space_used.sh
