# .config/systemd/user/systemd-analyze_verify.service
[Service]
ExecStart=%h/Programme/usw/shell/systemd-analyze_verify.sh
