# .config/systemd/user/unused/discharge_notification.timer
[Unit]
After = graphical.target

[Timer]
# Refresh the wallpaper every minute.
OnCalendar = *:0/01

[Install]
WantedBy = timers.target
