# .config/systemd/user/unused/net_logger.timer
[Unit]
After = network.target

[Timer]
OnCalendar = 00/0:01:00

[Install]
WantedBy = timers.target
