# .config/systemd/user/unused/net_logger.service
[Service]
ExecStart = %h/.local/bin/net_logger
