# .config/systemd/user/unused/wallpaperd.service
[Service]
ExecStart=/home/pactrag/Programme/usw/AWK/wallpaper.awk
# Do NOT omit this, it is important.
# Environment = DISPLAY = :0
